# ops/ci/jeryu-gc.service
[Unit]
Description=Jeryu SmartCache Garbage Collection
Documentation=https://github.com/jeppsontaylor/JeRyu
After=network.target docker.service

[Service]
Type=oneshot
User=ubuntu
ExecStart=/home/ubuntu/.local/bin/jeryu gc
TimeoutStartSec=600
# Log output to journal for observability
StandardOutput=journal
StandardError=journal
SyslogIdentifier=jeryu-gc
# Don't pollute the process table with failed runs
RemainAfterExit=no
